# docker/php/crontab
* * * * * root echo "Hello TDLAM" >> /var/log/cron.log 2>&1
